# vagrant/files/etc/profile.d/bash_aliases.sh
alias ls='ls --color=auto'
alias dir='ls -al'
alias ll='ls -l'
alias cd..='cd ..'
alias ..='cd ..'
